# build/package/debian/nginx-asg-sync.service
[Unit]
Description=NGINX Plus Integration with Cloud Autoscaling
After=nginx.service

[Service]
Type=simple
PermissionsStartOnly=true
User=nginx
Group=nginx
ExecStartPre=/bin/mkdir -p /var/log/nginx-asg-sync
ExecStartPre=/bin/chown nginx:nginx /var/log/nginx-asg-sync
ExecStart=/usr/sbin/nginx-asg-sync -log_path=/var/log/nginx-asg-sync/nginx-asg-sync.log
ExecStop=/bin/kill $MAINPID

[Install]
WantedBy=multi-user.target
